# com.kicktipp-ai.bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.kicktipp-ai.bot</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>cd /Users/hug/projects/kicktipp-mcp &amp;&amp; ./run.sh</string>
    </array>

    <!-- Every 6 hours. World Cup 2026 has multiple matchdays per week;
         6h catches everything with a comfortable buffer. -->
    <key>StartInterval</key>
    <integer>21600</integer>

    <key>RunAtLoad</key>
    <true/>

    <!-- launchd's bash -lc only sources /etc/profile, not ~/.zshrc, so
         ~/.local/bin (where `claude` lives) is missing from PATH. Point
         strategies.py's CLAUDE_CLI override at the absolute path instead
         of relying on PATH lookup. -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>CLAUDE_CLI</key>
        <string>/Users/hug/.local/bin/claude</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/hug/Library/Logs/kicktipp-ai.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/hug/Library/Logs/kicktipp-ai.log</string>
</dict>
</plist>
